# deploy/crawler-video.service
[Unit]
Description=Crawler FPV GStreamer Stream
After=network-online.target

[Service]
Type=simple
# No 'User=pi' needed here; it inherently runs as the user who owns it
ExecStart=/usr/bin/gst-launch-1.0 rpicamsrc bitrate=1000000 sensor-mode=7 keyframe-interval=15 preview=false inline-headers=true \
   ! 'video/x-h264,width=640,height=480,framerate=48/1,profile=baseline' \
   ! h264parse \
   ! rtph264pay config-interval=1 pt=96 mtu=1200 \
   ! udpsink host=192.168.1.111 port=5000 sync=false async=false

# Restart only if it crashes while it's SUPPOSED to be running
Restart=on-failure
RestartSec=2s

# Standard signals for a clean GStreamer exit
KillMode=process
Signal=SIGINT

[Install]
# For user services, this is the correct boot target
WantedBy=default.target
